# fishmod/fishmod_bootstrap.conf
description "System Startup Bind Script"
author "Ghostypeeps"

start on started ui

script
    # TODO: mount --bind fishmod stuff to root
    mount --bind /usr/local/fishmod /fsh
    mount --bind /fsh/bin/fiish.sh /usr/bin/crosh
    mount --bind /fsh/etc/chrome_dev.conf /etc/chrome_dev.conf
end script