# fishmod/fishmod_user_session.conf
description "User Login Bind Script"
author "Ghostypeeps"

start on start-user-session

script
    # TODO: mount --bind fishmod user data
end script